# packageroot/etc/systemd/system/climblive.service.d/override.conf
[Service]
Environment="DB_USERNAME=climblive"
Environment="DB_PASSWORD=climblive"
Environment="DB_HOST=localhost"
Environment="DB_PORT=3306"
Environment="DB_DATABASE=climblive"